# etc/cron.d/cluster-manager
# m h dom mon dow user  command
* * * * * root /usr/sbin/cluster-manager noop 2>&1 | logger -t cluster-manager
*/30 * * * * root /usr/sbin/cluster-manager check_sync_transcodes 2>&1 | logger -t cluster-manager
