# etc/cron.d/prt-cluster-manager
# m h dom mon dow user  command
#DEBUG * * * * * root /opt/plex/cluster_manager >> /var/log/prt-cluster-manager.log 2>&1
* * * * * root /opt/plex/cluster_manager
