# etc/cron.d/setup-plex-libraries
# m h dom mon dow user  command
*/30 * * * * root /usr/sbin/setup-plex-libraries
