# etc/cron.d/update-plex-libraries
# m h dom mon dow user  command
#DEBUG   23 1 * * * plex /usr/sbin/update-plex-libraries >> /var/lib/plexmediaserver/Library/Logs/update-plex-libraries.log 2>&1
0 * * * * plex /usr/sbin/update-plex-libraries
